# ubuntu-feed.timer
[Unit]
Description=Ubuntu Installer RSS Feed Generator

[Timer]
OnUnitActiveSec=1d
Unit=ubuntu-feed.service

[Install]
WantedBy=timers.target
